# fafo-erpnext/build/supervisord.conf
; Supervisord config for all Frappe processes inside the single
; ERPNext + Agriculture container. Loaded from
; /etc/supervisor/conf.d/erpnext.conf by the default supervisord
; that's already installed via apt in the Dockerfile.
;
; All processes:
;   * run as the `frappe` user (the base image's non-root user),
;   * autostart with the container,
;   * autorestart on any exit (Umbrel restarts the whole container
;     on failure via `restart: on-failure` in the compose, but
;     supervisord also restarts individual processes if they crash
;     without killing the whole container — faster recovery),
;   * log to stdout/stderr so `docker logs` shows everything.
;
; The `redis` program binds 127.0.0.1 only, so no external port
; needed — all inter-process calls go through localhost.

; ── Asset watchdog ─────────────────────────────────────────────────
; Polls sites/assets/ every 30 seconds. Restores from image snapshot
; when the canary file disappears — closes the gap where a runtime
; `bench migrate` / `bench build` / `bench clear-website-cache` wipes
; assets mid-container-lifetime, leaving nginx to serve unstyled
; pages until manual intervention.
[program:asset-watchdog]
command=/usr/local/bin/asset-watchdog
user=root
autostart=true
autorestart=true
startretries=99
priority=5
environment=WATCHDOG_INTERVAL="30"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Frappe Python backend via gunicorn. Priority 20 so it starts
; AFTER redis (which the backend connects to at startup). Two
; workers is a Pi-friendly default; bump if the host has more RAM.
[program:backend]
command=/home/frappe/frappe-bench/env/bin/gunicorn --bind 0.0.0.0:8000 --workers 2 --worker-class=gthread --threads 4 --timeout 300 --preload frappe.app:application
directory=/home/frappe/frappe-bench/sites
user=frappe
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Node socketio for real-time UI (typing indicators, live updates).
; Not strictly required but ERPNext UX degrades without it.
[program:websocket]
command=/home/frappe/.nvm/versions/node/v20.19.2/bin/node /home/frappe/frappe-bench/apps/frappe/socketio.js
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
priority=20
environment=NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Nginx frontend. Serves static assets directly, proxies /api/*
; to the backend gunicorn, proxies /socket.io/* to the node
; websocket. Copies BACKEND / SOCKETIO env vars into the nginx
; config via the shipping `nginx-entrypoint.sh`.
[program:frontend]
command=/usr/local/bin/nginx-entrypoint.sh
; Nginx frontend needs to run as root so it can open its log FDs
; and bind sub-1024 ports if the site config requests it. The nginx
; master process drops privileges to its own worker user (`nginx`
; or `nobody`) per its internal config after startup.
user=root
autostart=true
autorestart=true
priority=30
; FRAPPE_SITE_NAME_HEADER pinned to `frontend` (the site name we
; created via `bench new-site`) instead of `$host`. Frappe uses this
; to pick which site to serve; with `$host` any request whose Host
; header isn't `frontend` (like `umbrel.local:5300`) 404s because
; no site by that name exists. Pinning to `frontend` makes the
; container single-tenant — all requests get routed to that one
; site. If we ever go multi-tenant, revert to `$host` and create
; per-hostname sites.
environment=BACKEND="127.0.0.1:8000",FRAPPE_SITE_NAME_HEADER="frontend",SOCKETIO="127.0.0.1:9000",UPSTREAM_REAL_IP_ADDRESS="127.0.0.1",UPSTREAM_REAL_IP_HEADER="X-Forwarded-For",UPSTREAM_REAL_IP_RECURSIVE="off",PROXY_READ_TIMEOUT="300",CLIENT_MAX_BODY_SIZE="50m"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Background queue workers. Frappe's `bench worker` pulls jobs
; from redis DB 1 and executes them. Split into short + long
; queues per Frappe convention so short-lived jobs don't get
; stuck behind hour-long imports.
[program:queue-short]
command=/usr/local/bin/bench worker --queue short,default
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-long]
command=/usr/local/bin/bench worker --queue long
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Cron-like scheduler. Fires DocType events (daily / hourly / cron
; scheduled reports, background integrations) into the queue.
; Without this, scheduled jobs never run.
[program:scheduler]
command=/usr/local/bin/bench schedule
directory=/home/frappe/frappe-bench
user=frappe
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
